# src/fty-info.service.in
[Unit]
Description=fty-info service: serve metrics about this system as rackcontroller-0
Requires=malamute.service network.target
After=malamute.service network.target
Conflicts=rescue.target shutdown.target poweroff.target halt.target reboot.target emergency.target

Requisite=bios-allowed.target
After=bios-allowed.target
PartOf=bios-allowed.target
BindsTo=bios-allowed.target
Conflicts=bios-shutdown.target

Before=fty-asset.service
###UNDO the experiment###
### After=fty-asset.service
### Requires=fty-asset.service
### BindsTo=fty-asset.service

# fty-info/fty-mdns-sd should be pre-EULA, so that just-deployed instances
# are visible on the network for remote (distributed, through central point)
# commissioning, eventually.
PartOf=bios-pre-eula.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/fty/etc/default/fty
EnvironmentFile=-@prefix@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/fty
EnvironmentFile=-@sysconfdir@/default/fty__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-info -c @sysconfdir@/@PACKAGE@/fty-info.cfg

[Install]
WantedBy=bios-pre-eula.target
